# hdmi-autodetect/85-hdmi.rules
#! /bin/sh

ACTION=="change", TAG+="systemd", RUN+="/bin/systemctl start monitors.service" 
