# optionals/hdmi-autodetect/monitors.service
[Service]
Type=oneshot
ExecStart=/home/henrik/bin/hdmi-script
